# Ubuntu14/ansible-roles/dbielinski.confluent/files/schema-registry.conf
# schema-registry - schema-registry service job file
description "Schema-Registry Upstart Script"
author "Dennis Bielinski"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
# Automatically restart process if crashed
respawn
exec start-stop-daemon --start --chuid kafka --startas /bin/bash --make-pidfile --pidfile /var/log/kafka/schema-registry.pid -- -c "exec /usr/bin/schema-registry-start /etc/schema-registry/schema-registry.properties"
